# container-files/etc/supervisor.d/apache.conf
[program:apache]
command = /usr/sbin/httpd -DFOREGROUND
autorestart = true
